# scripts/supervisor.conf
[program:rtfd]
command=rtfd api
;user=
autostart=true
startsecs=5
exitcodes=0
redirect_stderr=true
