# scripts/supervisor.conf
[program:rtfd]
# If you start as a normal user, please fill in user
;user=
command=rtfd api
autostart=true
startsecs=5
exitcodes=0
redirect_stderr=true
